# data/systemd/gnome-flashback.service.in
[Unit]
Description=GNOME Flashback

After=gnome-session-manager.target

Requisite=gnome-session-initialized.target
PartOf=gnome-session-initialized.target
Before=gnome-session-initialized.target

StartLimitIntervalSec=15s
StartLimitBurst=3

[Service]
Type=dbus
Slice=session.slice
ExecStart=@bindir@/gnome-flashback
Restart=on-failure
RestartSec=0ms
BusName=org.gnome.Shell
