# data/autostart/gnome-flashback-screensaver.desktop.in
[Desktop Entry]
Type=Application
Name=Screensaver (GNOME Flashback)
Comment=Launch screensaver and locker program
# TRANSLATORS: Don't translate this text (this is icon name)
Icon=preferences-desktop-screensaver
Exec=gnome-screensaver
TryExec=gnome-screensaver
OnlyShowIn=GNOME-Flashback;
NoDisplay=true
X-GNOME-Autostart-Phase=Application
X-GNOME-Autostart-Notify=true
