# data/autostart/gnome-flashback-nm-applet.desktop.in
[Desktop Entry]
Name=Network (GNOME Flashback)
Comment=Manage your network connections
# TRANSLATORS: Don't translate this text (this is icon name)
Icon=nm-device-wireless
Exec=nm-applet
TryExec=nm-applet
Terminal=false
Type=Application
NoDisplay=true
OnlyShowIn=GNOME-Flashback;
X-GNOME-UsesNotifications=true
